# devnet/tasks.ini
[include]
files = %(here)s/cronos_777-1/tasks.ini

[supervisord]
pidfile = %(here)s/supervisord.pid
nodaemon = true
logfile = /dev/null
logfile_maxbytes = 0

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = %(here)s/supervisor.sock

[supervisorctl]
serverurl = unix://%(here)s/supervisor.sock

[program:relayer-demo]
autostart = false
autorestart = true
redirect_stderr = true
startsecs = 3
command = hermes -c %(here)s/relayer.toml start
stdout_logfile = %(here)s/relayer-demo.log

